# deployments/deb/distributed-encoder-agent.service
[Unit]
Description=Distributed Encoder Agent
Documentation=https://github.com/badskater/distributed-encoder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=distributed-encoder-agent
Group=distributed-encoder-agent

# Optional: place DE_* environment variables here instead of (or alongside)
# the YAML config file. The leading '-' makes this file optional.
EnvironmentFile=-/etc/distributed-encoder/agent-environment

ExecStart=/usr/bin/distributed-encoder-agent run \
  --config /etc/distributed-encoder/agent.yaml

Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=distributed-encoder-agent

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ReadWritePaths=/var/lib/distributed-encoder-agent /var/log/distributed-encoder-agent

[Install]
WantedBy=multi-user.target
